# docker/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
logfile=/home/devops/logs/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[include]
files = /etc/supervisor/conf.d/*.conf

[program:bench]
command=/bin/bash -c "cd /home/devops/frappe-bench && bench start"
directory=/home/devops/frappe-bench
autostart=true
autorestart=true
stdout_logfile=/home/devops/logs/bench.log
stderr_logfile=/home/devops/logs/bench.err.log

